# crontab
30 4 * * * cd /home/pi/Documents/KismetUIPlugin ; ./visualise.sh /home/pi/Documents/Share/wifi ; rm probes.csv ; ./start.sh
