# crontab
@reboot sleep 300 && cd /home/pi/Documents/KismetUIPlugin && ./start.sh
5 0 * * * cd /home/pi/Documents/KismetUIPlugin && ./visualise.sh /home/pi/Documents/Share/wifi

